# deploy/dewdrop-backup.timer
[Unit]
Description=Run DEWDROP SQLite backup nightly

[Timer]
# Nightly at 01:00 local time (after scoring + ensemble snapshot); catch up if off.
OnCalendar=*-*-* 01:00:00 America/Chicago
Persistent=true

[Install]
WantedBy=timers.target
