# deploy/dewdrop-backup.service
[Unit]
Description=DEWDROP nightly SQLite backup
After=network-online.target

[Service]
Type=oneshot
User=dewdrop
WorkingDirectory=/opt/dewdrop
ExecStart=/opt/dewdrop/venv/bin/python scripts/backup_db.py
EnvironmentFile=/opt/dewdrop/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dewdrop-backup

# ── Sandboxing / hardening ────────────────────────────────────────────────
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
# If DEWDROP_BACKUP_DIR points outside /opt/dewdrop/data, add it here too.
ReadWritePaths=/opt/dewdrop/data /opt/dewdrop/logs
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
UMask=0077
