# deploy/dewdrop-actuals.service
[Unit]
Description=DEWDROP actuals ingester (ASOS/MCI + EcoWitt, design doc §4 Phase 2)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=dewdrop
WorkingDirectory=/opt/dewdrop
ExecStart=/opt/dewdrop/venv/bin/python scripts/ingest_actuals.py
EnvironmentFile=/opt/dewdrop/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dewdrop-actuals

# ── Sandboxing / hardening ────────────────────────────────────────────────
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/dewdrop/data /opt/dewdrop/logs
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
UMask=0077
