# deploy/dewdrop-station.timer
[Unit]
Description=Poll GW2000X weather station every minute

[Timer]
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target
