# deploy/dewdrop-score.timer
[Unit]
Description=Run DEWDROP forecast scoring nightly

[Timer]
# Nightly at 00:45 local time (after actuals at 00:30); catch up if off.
OnCalendar=*-*-* 00:45:00 America/Chicago
Persistent=true

[Install]
WantedBy=timers.target
