# deploy/dewdrop-actuals.timer
[Unit]
Description=Fetch yesterday's actuals nightly (design doc §4 Phase 2)

[Timer]
# After midnight, once the previous day is complete at all sources.
OnCalendar=*-*-* 00:30:00 America/Chicago
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
